# dracut/28immucore/immucore
#!/sbin/openrc-run
  
name=$RC_SVCNAME
command="/usr/bin/immucore"

depend() {
        need udev
}
