# dracut/28immucore/kairos-setup-initramfs
#!/sbin/openrc-run

name=$RC_SVCNAME
description="kairos system initramfs setup before switch root"
command="/usr/bin/elemental"
command_args="run-stage initramfs"
start_stop_daemon_args=" -r /sysroot"

depend() { 
	need root 
}
